# examples/io.lassi.privatetcp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>io.lassi.privatetcp</string>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/local/libexec/privatetcp</string>
                <string>http</string>
        </array>
        <key>KeepAlive</key>
        <true/>
        <key>RunAtLoad</key>
        <true/>
        <key>UserName</key>
        <string>root</string>
        <key>StandardErrorPath</key>
        <string>/var/log/privatetcp.log</string>
</dict>
</plist>
